# salt/httpd/tmpwatch-results.cron

 0  * * * *  root  tmpwatch 10d /srv/web/results/ >/var/log/tmpwatch-results.log 2>&1
 10 * * * *  root  tmpwatch 5d /home/gw/airavata-workdirs/ >/var/log/tmpwatch-airavata.log 2>&1
 20 * * * *  root  tmpwatch 10d /var/lib/ganglia/rrds/ >/var/log/tmpwatch-ganglia.log 2>&1

