# salt/tools/sra-report-to-es.cron

*/5 * * * *  root  /srv/jetstream-sra-cluster/salt/tools/sra-report-to-es >/dev/null 2>&1

